# scripts/provider-remote-https/oasis7-remote-provider-bridge.service
[Unit]
Description=oasis7 remote HTTPS provider bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=oasis7-provider
Group=oasis7-provider
WorkingDirectory=/opt/oasis7/remote-provider-bridge/current
EnvironmentFile=/etc/oasis7/remote-provider-bridge.env
ExecStart=/opt/oasis7/remote-provider-bridge/current/scripts/provider-remote-https/start-remote-provider-bridge.sh
Restart=always
RestartSec=3
TimeoutStopSec=15
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
MemoryDenyWriteExecute=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
